# scripts/vllm-fwd.service
[Unit]
Description=vLLM WSL socket forwarder (0.0.0.0:18000 -> k8s NodePort)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rjmendez
# Resolve the k8s NodePort target at EACH start. Live-verified 2026-07-08:
# the NodePort (32272) is served ONLY on the primary /20 eth0 address
# (172.21.171.198), NOT on 127.0.0.1, localhost, 10.255.255.254, or the
# static /32 172.31.255.254. That primary IP is DHCP-assigned and can change
# on WSL restart, so we derive it dynamically instead of hardcoding.
Environment=VLLM_FWD_PORT=18000
Environment=VLLM_FWD_TARGET_PORT=32272
ExecStart=/bin/bash -c 'export VLLM_FWD_TARGET_HOST=$(ip route get 1.1.1.1 | awk "{print \$7; exit}"); echo "target=$VLLM_FWD_TARGET_HOST:$VLLM_FWD_TARGET_PORT"; exec /usr/bin/python3 /home/rjmendez/development/loci/scripts/vllm_tailscale_forward.py'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
